# output/windows/process_creation/scheduled_task_creation_on_remote_endpoint_using_at.yml
title: Scheduled Task Creation on Remote Endpoint using At
id: 4be54858-432f-11ec-8209-3e22fbd008af
status: stable
level: medium
author: Mauricio Velazco, Splunk
description: The following analytic detects the creation of scheduled tasks on remote Windows endpoints using the at.exe command. This detection leverages Endpoint Detection and Response (EDR) telemetry,
  focusing on process creation events involving at.exe with remote command-line arguments. Identifying this activity is significant for a SOC as it may indicate lateral movement or remote code execution
  attempts by an attacker. If confirmed malicious, this activity could lead to unauthorized access, persistence, or execution of malicious code, potentially resulting in data theft or further compromise
  of the network.
references:
- https://docs.microsoft.com/en-us/windows-server/administration/windows-commands/at
- https://docs.microsoft.com/en-us/windows/win32/cimwin32prov/win32-scheduledjob?redirectedfrom=MSDN
tags:
- attack.T1053_002
date: '2021-11-11'
logsource:
  category: process_creation
  product: windows
detection:
  selection_0:
    Image: at.exe
  selection_1:
    OriginalFileName: at.exe
  selection_2:
    CommandLine|contains: '*\\*'
  condition: selection_0 and selection_1 and selection_2
