# output/windows/process_creation/suspicious_scheduled_task_from_public_directory.yml
title: Suspicious Scheduled Task from Public Directory
id: 7feb7972-7ac3-11eb-bac8-acde48001122
status: stable
level: medium
author: Michael Haag, Splunk
description: The following analytic identifies the creation of scheduled tasks that execute binaries or scripts from public directories, such as users\public, \programdata\, or \windows\temp, using schtasks.exe
  with the /create command. It leverages Sysmon Event ID 1 data to detect this behavior. This activity is significant because it often indicates an attempt to maintain persistence or execute malicious scripts,
  which are common tactics in malware deployment. If confirmed as malicious, this could lead to data compromise, unauthorized access, and potential lateral movement within the network.
references:
- https://attack.mitre.org/techniques/T1053/005/
tags:
- attack.T1053_005
date: '2021-03-01'
logsource:
  category: process_creation
  product: windows
detection:
  selection_0:
    Image: schtasks.exe
  selection_1:
    CommandLine|contains: '*\users\public\*'
  selection_2:
    CommandLine|contains: '*\programdata\*'
  selection_3:
    CommandLine|contains: '*windows\temp*'
  selection_4:
    CommandLine|contains: '*/create*'
  condition: selection_0 and selection_1 and selection_2 and selection_3 and selection_4
